# system/99-capslock-led.rules
ACTION=="add", SUBSYSTEM=="leds", KERNEL=="*::capslock", RUN+="/bin/chmod 0666 %S%p/brightness", RUN+="/bin/sh -c 'echo none > %S%p/trigger'"
